# examples/openclaw-search-bridge.service.example
[Unit]
Description=OpenClaw Search Bridge
After=network.target

[Service]
Type=simple
Environment=BRAVE_API_KEY=YOUR_BRAVE_API_KEY
Environment=HTTP_PROXY=http://172.27.208.1:10808
Environment=HTTPS_PROXY=http://172.27.208.1:10808
Environment=NO_PROXY=127.0.0.1,localhost,::1
ExecStart=/usr/bin/python3 /home/<your-user>/.openclaw/search-bridge-mvp/search_bridge.py
WorkingDirectory=/home/<your-user>/.openclaw/search-bridge-mvp
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
